# zfs_uploader.service
[Unit]
Description=ZFS snapshot to S3 uploader
Wants=network-online.target

[Service]
Environment=PYTHONUNBUFFERED=1
ExecStart=/etc/zfs_uploader/env/bin/zfsup backup
Restart=on-failure
WorkingDirectory=/etc/zfs_uploader

[Install]
WantedBy=default.target
